# debian/postinst
#!/bin/sh -ex

systemctl enable elasticsearch
systemctl start elasticsearch
