# debian/postinst
#!/bin/sh -ex

update-rc.d elasticsearch defaults
/etc/init.d/elasticsearch start
